# scripts/supervisord.conf
[supervisord]
nodaemon=true

[program:masarService]
command=/build/masarService/cpp/bin/linux-x86_64/masarServiceRun masarService
